# workloads/container-shared.service
[Unit]
StopWhenUnneeded=yes

[Service]
ExecStart=/bin/sleep 1000d