# workloads/workload-c@node1.target
# Workload C on node 1
[Unit]
StopWhenUnneeded=yes

# Services in this workload
Requires=container-c1.service