# workloads/workload-a@node1.target
# Workload A on node 1
[Unit]
StopWhenUnneeded=yes

# Services in this workload
Requires=container-a1.service
Requires=container-shared.service