# workloads/workload-c@.target
# Default workload C on generic nodes (no deps)
[Unit]
StopWhenUnneeded=yes

# No services in this workload