# workloads/workload-c@node2.target
# Workload C on node 2
[Unit]
StopWhenUnneeded=yes

# Services in this workload
Requires=container-c2.service